# xremaprunner.service
[Unit]
Description=Xremap Runner

[Service]
ExecStart=/usr/bin/xrmp-runner
Type=simple
WorkingDirectory=/
Restart=no