# res/init/openrc/xremaprunner.in
#!/sbin/openrc-run

name="xremaprunner"
description="Xremap Runner"
command=@bindir@/xrmp-runner
command_background=true
