# res/init/systemd/xremaprunner.service.in
[Unit]
Description=Xremap Runner

[Service]
ExecStart=@bindir@/xrmp-runner
Type=simple
WorkingDirectory=/
Restart=no
